# yes/etc/profile.d/conda.csh
setenv CONDA_EXE "/home/iven/igb-design-center/yes/bin/conda";
setenv _CONDA_EXE "/home/iven/igb-design-center/yes/bin/conda";
unsetenv _CE_M;
unsetenv _CE_CONDA;
setenv CONDA_PYTHON_EXE "/home/iven/igb-design-center/yes/bin/python";
setenv _CONDA_ROOT "/home/iven/igb-design-center/yes";

source "/home/iven/igb-design-center/yes/lib/python3.13/site-packages/conda/shell/etc/profile.d/conda.csh"